# scripts/photonic-backend.service
[Unit]
Description=Photonic Photobooth Backend
After=network.target

[Service]
Type=simple
User=qiu
WorkingDirectory=/home/qiu/photonic-v0.1/apps/backend
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
